# qr3_33.conf
[program:qr3_33]
command=python3 /home/curator3/qr3_33/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_33/test.err.log
stdout_logfile=/home/curator3/qr3_33/test.err.log
